# crontab
10-59/15 * * * * /usr/local/bin/hs -c "hs.alert.show('drink some water')"
