# crontab
10-59/15 * * * * osascript -e 'display notification "drink some water" with title "Health Tips"'
